# deploy/systemd/vectorcore-sgw-c.service
[Unit]
Description=VectorCore SGW-C
Documentation=/opt/vectorcore/docs/production-controls.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/vectorcore/bin/sgw-c -c /opt/vectorcore/etc/sgw-c.yaml
Restart=on-failure
RestartSec=5s
User=root
Group=root
LimitNOFILE=1048576
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/vectorcore-sgw /var/log/vectorcore /opt/vectorcore/log

[Install]
WantedBy=multi-user.target
